# systemd/.config/systemd/user/repo-watch.path
[Path]
PathChanged=/home/axionl/Programs/archlinux/local/public

[Install]
WantedBy=default.target